# packages/tvtv/tmp/DEBIAN/postinst
#!/bin/sh

echo "creating init levels for Reel TvTv client daemon"
update-rc.d reeltvtv defaults

/etc/init.d/reeltvtv start

