# packages/wlan-rtl8192/files
#enter all files here that have to be copied from the SRCDIR to the package root
#note: use only tabs and only one tab between fields!
#
#valid types:
# script - for scripts, default 755, root.root
# exe    - like script, but will be stripped on packaging
# sexe   - setuid root executeable, 4755, root.root, strip
# lib	 - library,     755, root.root, no strip
# cfg    - config file, 644, root.root
# bin    - binaries, i.e. firmware, 644, root.root
# any other files will be installed "as they are", so be careful about permissions!
#
# src rel. to $(SRCDIR)	dest			type	owner	group

#etc/init.d/skeleton	/etc/init.d/skeleton	script	root	root
# this would be a more simple form:
# etc/init.d/skeleton	-			script

