# packages/djmount/postrm
#!/bin/sh

if [ "$1" = "purge" ] ; then
	update-rc.d djmount remove >/dev/null
fi

exit 0
