# packages/ismd/ROOT/etc/udev/rules.d/81-ice-smd.rules
# do not edit this file, it will be overwritten on update

ACTION!="add", GOTO="ice_smd_drivers_end"

SUBSYSTEM=="module", KERNEL=="clock_control", RUN+="/sbin/create_dev.sh make_dev clock_control clock_control"
SUBSYSTEM=="module", KERNEL=="ioctl_module", RUN+="/sbin/create_dev.sh make_dev smdmsg smdmsg", RUN+="chmod 666 /dev/smdmsg"

LABEL="ice_smd_drivers_end"
