# packages/mediatomb/postinst
#!/bin/sh

chown mediatomb:mediatomb /etc/mediatomb/config.xml
chown mediatomb:mediatomb /var/lib/mediatomb

# Automatically added by dh_installinit
if [ -x "/etc/init.d/mediatomb" ]; then
	update-rc.d mediatomb defaults >/dev/null
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d mediatomb start
	else
		/etc/init.d/mediatomb start
	fi
fi
# End automatically added section


exit 0
