# packages/usbmount/tmp/etc/udev/rules.d/80-ptp-camera-mount.rules
# Rules for PTP camera
PROGRAM="check-ptp-camera 06/01/01", ACTION=="add", SUBSYSTEM=="usb_device", RUN+="/usr/sbin/PTPCameraMount.sh add"
ACTION=="remove", SUBSYSTEM=="usb_device", RUN+="/usr/sbin/PTPCameraMount.sh remove"
