# packages/reel-vdradmin/tmp/etc/cron.d/vdradmin-am
MAILTO=""
30 6 * * * root /etc/init.d/vdradmin restart
