# packages/reelbox/postrm
#! /bin/sh

if [ "$1" = "purge" ]
then
    update-rc.d sanitycheck-reel remove >/dev/null
fi

exit 0
