# packages/bftpd/postinst
#!/bin/sh

echo "creating init levels for btftp server daemon"
update-rc.d bftpd-server defaults 04 80

/etc/init.d/bftpd-server start

