# packages/display/tmp/etc/udev/rules.d/81-ice-display.rules
# do not edit this file, it will be overwritten on update

ACTION!="add", GOTO="ice_display_drivers_end"

SUBSYSTEM=="module", KERNEL=="avcap_core", RUN+="/sbin/create_dev.sh make_dev avcap_core avcap_core"
SUBSYSTEM=="module", KERNEL=="gdl_server", RUN+="/sbin/gdl-server-starter.sh"
#SUBSYSTEM=="module", KERNEL=="sec_kernel", RUN+="/sbin/create_dev.sh make_dev sec sec", RUN+="/sbin/load_hdmi.sh"
SUBSYSTEM=="module", KERNEL=="sec_kernel", RUN+="/sbin/create_dev.sh make_dev sec sec", RUN+="/sbin/load_sec.sh"
#obs in pc21 SUBSYSTEM=="module", KERNEL=="pvrsrvkm", RUN+="/sbin/create_dev.sh make_dev pvrsrvkm pvrsrvkm", RUN+="pvrsrvinit"
SUBSYSTEM=="module", KERNEL=="pvrsrvkm", RUN+="/usr/sbin/pvrsrvinit"

LABEL="ice_display_drivers_end"
