# packages/platform-config/tmp/etc/udev/rules.d/81-ice-platform-config.rules
# do not edit this file, it will be overwritten on update

ACTION!="add", GOTO="ice_drivers_end"

#DRIVER!="?*", ENV{MODALIAS}=="?*", RUN+="/sbin/modprobe -b $env{MODALIAS}"
SUBSYSTEM=="module", KERNEL=="platform_config", RUN+="/sbin/platform_config.sh start"

LABEL="ice_drivers_end"
