# packages/telnetd-reel/postrm
#! /bin/sh

set -e

if [ "$1" = "purge" ]; then
    update-rc.d telnetd-reel remove >/dev/null
fi

