# packages/dyndns-client/postrm
#! /bin/sh

set -e

if [ "$1" = "purge" ]; then
    update-rc.d reel-dyndns remove >/dev/null
fi

