# packages/reel-wlan/tmp/DEBIAN/postrm
#! /bin/sh

update-rc.d -f reel-wlan remove >/dev/null

exit 0
