# packages/system-utils/ROOT/etc/udev/rules.d/81-ice-utils.rules
# do not edit this file, it will be overwritten on update

ACTION!="add", GOTO="ice_utils_drivers_end"

SUBSYSTEM=="module", KERNEL=="edl_thermal", RUN+="/sbin/create_dev.sh make_dev thermal thermal"
SUBSYSTEM=="module", KERNEL=="intel_ce_flash", RUN+="/sbin/create_dev.sh make_dev intel_ce_flash flash"
SUBSYSTEM=="module", KERNEL=="intel_ce_pm", RUN+="/sbin/create_dev.sh make_dev intel_ce_pm intel_ce_pm"
SUBSYSTEM=="module", KERNEL=="expbus", RUN+="/sbin/create_dev.sh make_dev expbus expbus"
###TODO: check if this works correct
SUBSYSTEM=="module", KERNEL=="iosf", RUN+="/sbin/create_dev.sh make_dev_minor iosf iosf0 0"

LABEL="ice_utils_drivers_end"
