# packages/tftp-hpa/tmp/DEBIAN/postinst
#!/bin/sh

echo "creating init levels for tftp server daemon"
update-rc.d tftpd defaults 04 80

/etc/init.d/tftpd start

