# packages/reel-umts/tmp/DEBIAN/postrm
#! /bin/sh

update-rc.d reel-umts remove >/dev/null

exit 0
