# packages/usb-display-firmware/tmp/etc/udev/rules.d/40-Reel_OLED.rules
# UDEV-Rule for Reel USB OLED display
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="054b", SYMLINK+="reel_oled", RUN+="/sbin/reel-usb-display-updater"
