# packages/reelbox-frontpanel-mini/ROOT/etc/udev/rules.d/50-reel_frontpanel.rules
# Rules for ReelBox frontpanel

KERNEL=="ttyS1", SUBSYSTEM=="tty", ACTION=="add",    RUN+="/sbin/dev_frontpanel.sh"
