# packages/idl/tmp/etc/udev/rules.d/81-idl.rules
# do not edit this file, it will be overwritten on update

ACTION!="add", GOTO="idl_drivers_end"

# DEPRECATED in PC25
SUBSYSTEM=="module", KERNEL=="idl_i2c", RUN+="/sbin/create_dev.sh make_dev i2c i2c"
# DEPRECATED in PC25
SUBSYSTEM=="module", KERNEL=="idl_spi", RUN+="/sbin/create_dev.sh make_dev spi spi0"
SUBSYSTEM=="module", KERNEL=="idl_gpio", RUN+="/sbin/create_dev.sh make_dev gpio gpio", RUN+="/bin/chmod 666 /dev/gpio"
SUBSYSTEM=="module", KERNEL=="idl_smartcard", RUN+="/sbin/create_dev.sh make_dev_minor scard0 scard0"

LABEL="idl_drivers_end"
