# packages/reelbox-ice/postrm
#! /bin/sh

if [ "$1" = "purge" ]
then
    update-rc.d sanitycheck-reel remove >/dev/null
fi

# installation scripts should always exit with 0. Otherwise it will create ugly
#  kinds of problems for the end-user.

exit 0
