# packages/netceiver-client/postrm
#! /bin/sh

if [ "$1" = "purge" ]; then
    update-rc.d -f mcli remove >/dev/null
fi

sed -i /dvbloop/d /etc/modules
