# packages/pal-osal/ROOT/etc/udev/rules.d/81-ice-osal.rules
# do not edit this file, it will be overwritten on update

ACTION!="add", GOTO="osal_drivers_end"

SUBSYSTEM=="module", KERNEL=="osal_linux", RUN+="/sbin/create_dev.sh make_dev devmem devmem", RUN+="chmod 666 /dev/devmem"

LABEL="osal_drivers_end"
